# pi-setup/netatmo.cron
# Every 5 min 7am - 9pm

*/5 07-21 * * * cd /home/pi/netatmo-pi/ && ./update-screen.sh
